# packaging/aleph-vm/DEBIAN/postinst
#!/bin/bash
set -euf -o pipefail

if ! id -u jailman > /dev/null 2>&1; then
  useradd jailman
fi

rm -fr /srv/jailer  # Upgrade from < 0.1.11
rm -fr /tmp/aleph   # Upgrade from < 0.1.11
mkdir -p /var/lib/aleph/vm/jailer

# Systemd is absent from containers
if ! [[ -v container ]]; then
  systemctl daemon-reload
  systemctl enable ipfs.service
  systemctl enable aleph-vm-supervisor.service
  systemctl restart aleph-vm-supervisor.service
fi
