# packaging/aleph-vm/etc/systemd/system/aleph-vm-controller@.service
[Unit]
Description=Aleph VM %i Controller
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 -m aleph.vm.controllers --config=/var/lib/aleph/vm/%i-controller.json
Restart=on-failure

[Install]
WantedBy=multi-user.target
